# scripts/library/launchd.uk.changsung.library-assistant-actions.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>uk.changsung.library-assistant-actions</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/curl</string>
    <string>-sS</string>
    <string>-X</string>
    <string>POST</string>
    <string>http://127.0.0.1:8802/api/library/apply</string>
    <string>-H</string>
    <string>Content-Type: application/json</string>
    <string>--data</string>
    <string>{"dry_run": false}</string>
  </array>

  <!-- Run twice daily: 12:30 and 00:40 -->
  <key>StartCalendarInterval</key>
  <array>
    <dict>
      <key>Hour</key><integer>12</integer>
      <key>Minute</key><integer>30</integer>
    </dict>
    <dict>
      <key>Hour</key><integer>0</integer>
      <key>Minute</key><integer>40</integer>
    </dict>
  </array>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/tmp/library-assistant-actions.out.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/library-assistant-actions.err.log</string>
</dict>
</plist>
